# systemd/telegram-claude-bot.service
[Unit]
Description=Telegram Claude Code Proxy Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/projects/telegram-claude-bot
Environment="TELEGRAM_BOT_TOKEN=YOUR_BOT_TOKEN_HERE"
Environment="ALLOWED_USER_IDS=YOUR_TELEGRAM_ID_HERE"
Environment="PROJECT_DIR=/home/ubuntu/projects/chefvision-analytics-platform"
ExecStart=/usr/bin/python3 /home/ubuntu/projects/telegram-claude-bot/telegram_proxy.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-claude-bot

[Install]
WantedBy=multi-user.target
